# launchd/com.rin.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rin.proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>__RIN_HOME__/src/rin_proxy/rin-proxy</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>GEMINI_API_KEY</key>
        <string>__GEMINI_API_KEY__</string>
        <key>GLM_API_KEY</key>
        <string>__GLM_API_KEY__</string>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.rin/proxy.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.rin/proxy.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>
</dict>
</plist>
