# launchd/com.rin.session-harvest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rin.session-harvest</string>
    <key>ProgramArguments</key>
    <array>
        <string>__RIN_HOME__/.venv/bin/python</string>
        <string>__RIN_HOME__/scripts/session-harvest.py</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/rin-session-harvest.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/rin-session-harvest.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>WorkingDirectory</key>
    <string>__RIN_HOME__</string>
</dict>
</plist>
